# etc/udev/rules.d/11-sd.rules
KERNEL=="sd[a-z]", RUN+="/sbin/udevh-sd %k"
KERNEL=="sr[0-9]", RUN+="/sbin/udevh-sd %k"
